# ansible/linux/playbook/deploy/deploy_node_exporter/deploy_node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=swnodeexp
Group=swnodeexp
type=simple
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=multi-user.target